# scripts/systemd/npm-log-worker.service
[Unit]
Description=NPM Monitor Log Worker (Real-time)
After=network.target postgresql.service redis.service

[Service]
Type=simple
User=mattia
WorkingDirectory=/home/mattia/npm-stats
EnvironmentFile=/home/mattia/npm-stats/.env
ExecStart=/home/mattia/npm-stats/.venv/bin/python3 -m src.log_worker
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
